# db/refresh_events.cron
0,15,30,45 * * * * /home/mnutt/limedirectory/script/refresh_events
